# crates/flight-updater/tests/fixtures/installer/flightd.service
[Unit]
Description=Flight Hub daemon
Documentation=https://github.com/EffortlessMetrics/OpenFlight

[Service]
Type=simple
ExecStart=/usr/bin/flightd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
